# supervisorctl.ini
[program:newpinvis]
command=bash -c 'echo "waiting 28 days ..." && sleep 28d && exec /path/to/your/rebuild_site.sh'
autostart=yes
autorestart=yes
startsecs=5
stopasgroup=true

